# ops/systemd/signaris-hub-course-due-soon.timer
[Unit]
Description=Run signaris-hub-course-due-soon daily 09:15 Moscow time

[Timer]
# 06:15 UTC = 09:15 MSK (UTC+3, no DST in RU)
OnCalendar=*-*-* 06:15:00
Persistent=true
Unit=signaris-hub-course-due-soon.service

[Install]
WantedBy=timers.target
