# ops/systemd/signaris-hub-staging-automations.service
[Unit]
Description=Signaris Hub — learn automations runner (-staging hourly cron)
After=network.target postgresql.service

[Service]
Type=oneshot
User=signaris
Group=signaris
WorkingDirectory=/opt/signaris-hub-staging
EnvironmentFile=/opt/signaris-hub-staging/.env
ExecStart=/opt/signaris-hub-staging/.venv/bin/python -m app.jobs.automations_run
